# server/devpi_server/cfg/supervisord.conf.template

[program:devpi-server]
command={devpibin} {server_args}
user = {user}
priority = 999
startsecs = 5
redirect_stderr = True
autostart = True
