# server/devpi_server/cfg/supervisor-devpi.conf.template
[program:devpi-server]
command={devpibin} {server_args}
priority = 999
startsecs = 5
redirect_stderr = True
autostart = True
